# build/raspbian/modules/configure/files/99-reterminal.rules
# Allow smart-panel service to control reTerminal LEDs and buzzer
SUBSYSTEM=="leds", KERNEL=="usr_led*", ACTION=="add", RUN+="/bin/chmod 0666 /sys%p/brightness"
SUBSYSTEM=="leds", KERNEL=="usr_buzzer", ACTION=="add", RUN+="/bin/chmod 0666 /sys%p/brightness"
